# bmc/meta-myhello/recipes-phosphor/hello/files/hello-1.0/hello.service
[Unit]
Description=Hello OpenBMC Meson Demo

[Service]
Type=simple
ExecStart=/usr/bin/hello
Restart=always

[Install]
WantedBy=multi-user.target